# launchd/com.user.singbox.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchDaemon для sing-box на macOS.
  Кладётся в /Library/LaunchDaemons/com.user.singbox.plist (root:wheel, 644).
  RunAtLoad + KeepAlive => поднимается при загрузке (до входа в систему) и
  сам перезапускается после сна/сбоя. Это и есть "открыл Mac -> уже работает".

  Демон запускает обёртку /etc/sing-box/daemon-run.sh, которая при включённом
  kill-switch поднимает защиту ещё до старта sing-box, затем запускает сам sing-box.
  install-macos-daemon.sh подставляет путь обёртки в __DAEMON_RUN__.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.singbox</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__DAEMON_RUN__</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Interactive</string>

    <key>StandardOutPath</key>
    <string>/var/log/sing-box.out.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/sing-box.err.log</string>
</dict>
</plist>
